# backend-express/rc.d/lorem.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: lorem
# REQUIRE: LOGIN

. /etc/rc.subr

name="lorem"
rcvar="lorem_enable"

pidfile="@APP_RUNDIR@/pid"
command="@prefix@/sbin/${name}"
command_args="--daemon"
procname="@node@"
#command_interpreter="@node@"

load_rc_config ${name}

: ${lorem_enable:="NO"}

run_rc_command "$1"
#EOF
