# backend/rc.d/wcm.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: wcm
# REQUIRE: LOGIN

. /etc/rc.subr

name="wcm"
rcvar="wcm_enable"

pidfile="@APP_RUNDIR@/pid"
command="@prefix@/sbin/${name}"
command_args="--daemon"
procname="@node@"
#command_interpreter="@node@"

load_rc_config ${name}

: ${wcm_enable:="NO"}

run_rc_command "$1"
#EOF
